# systemd/user/jarvis.service
[Unit]
Description=Jarvis Transcription Daemon
Requires=jarvis.socket
After=jarvis.socket

[Service]
# Path to your conda or system python
ExecStart=/home/mchang/miniconda3/bin/python3 /home/mchang/repos/nerd-dictation/scripts/jarvis_daemon.py
StandardOutput=journal
StandardError=journal

# Security Hardening
# PrivateTmp=yes # Disabled to allow sharing /tmp/jarvis_command.wav
# ProtectHome=yes # Disabled to allow reading from home dir if needed
ProtectSystem=full
NoNewPrivileges=yes

[Install]
WantedBy=default.target
