# systemd/user/jarvis.socket
[Unit]
Description=Jarvis Transcription Socket

[Socket]
ListenStream=%t/jarvis.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
